# pkg/linuwu-sense-gui/etc/udev/rules.d/60-linuwu-sense.rules
# 60-linuwu-sense.rules
# Grant the `acer-nitro` group read/write access to all linuwu_sense sysfs
# attributes so the GUI can run as a normal user with no root / pkexec.

# ── acer-wmi platform device (nitro_sense, predator_sense, four_zoned_kb) ──
ACTION=="add|change", SUBSYSTEM=="platform", KERNEL=="acer-wmi", \
    RUN+="/usr/lib/linuwu-sense-gui/set-sysfs-perms.sh %S%p"

# ── ACPI platform profile ─────────────────────────────────────────────────
ACTION=="add|change", SUBSYSTEM=="platform_profile", \
    RUN+="/bin/sh -c 'chgrp acer-nitro /sys/firmware/acpi/platform_profile /sys/firmware/acpi/platform_profile_choices 2>/dev/null; chmod 0664 /sys/firmware/acpi/platform_profile 2>/dev/null; chmod 0644 /sys/firmware/acpi/platform_profile_choices 2>/dev/null || true'"

# ── cpufreq scaling_max_freq (for CPU presets) ───────────────────────────
ACTION=="add|change", SUBSYSTEM=="cpu", KERNEL=="cpu[0-9]*", \
    RUN+="/bin/sh -c 'chgrp acer-nitro /sys/devices/system/cpu/cpufreq/policy%n/scaling_max_freq 2>/dev/null; chmod 0664 /sys/devices/system/cpu/cpufreq/policy%n/scaling_max_freq 2>/dev/null || true'"
